# deploy/games/dregg-games-bot.service
# ⚠⚠ ASPIRATIONAL — NOT INSTALLED ON hbox. ⚠⚠
#
# Verified 2026-07-15: hbox runs NO Discord bot. The bot runs ON THE AWS EDGE as the
# container `dreggnet-dreggnet-discord-bot-1` (image `dregg-discord-bot:staging`), from
# the compose stack at /opt/dreggnet — NOT from any systemd unit anywhere.
#
# The double-fire hazard noted below is REAL and still current, but its subject has
# moved: the token is held by the EDGE CONTAINER, not by a "graviton bot" systemd unit
# (deploy/aws/SUPERSEDED/dregg-discord-bot.service never ran). Before starting this unit:
#   ssh <edge> 'cd /opt/dreggnet && docker compose stop dreggnet-discord-bot'
#
# Moving the bot OFF the edge is TODO-2 in deploy/README.md — but persvati is the
# better target than hbox (hbox is already over-tenanted; see deploy/PRACTICES.md §1).
#
# ---------------------------------------------------------------------------
# dregg games Discord bot — the Descent daily + slash commands on hbox.
#
# Runs `dregg-discord-bot` (package `dregg-discord-bot`): the daily-reveal cron
# (live drand quicknet, BLS-verified), `/descent`, the semi-private channels, the
# no-cheat board. STANDALONE — the games bot does NOT require the dregg node (the
# demo verifies by replay, not by a node). Contrast deploy/aws/dregg-discord-bot.service,
# which `Requires=dregg-gateway.service`; here there is no node dependency.
#
# hbox discipline (deploy/hbox/RUNBOOK.md): USER unit, `%h` paths,
# `loginctl enable-linger hbox`, WantedBy=default.target. Mirrors the existing
# deploy/hbox/dregg-discord-bot.service, extended with the games env + hardening.
#
# ⚠ DOUBLE-FIRE HAZARD: two bots on ONE Discord token fire every command twice.
# The token is currently held by the GRAVITON bot (deploy/aws/dregg-discord-bot.service);
# hbox runs NO bot today. STOP the graviton bot FIRST before this one starts. See
# deploy/games/RUNBOOK.md step (c) — ember-gated. (This is a STANDALONE hbox unit with
# no gateway route — the bot talks straight to Discord + drand egress, not through Caddy.)
[Unit]
Description=dregg games Discord bot (Descent daily, hbox, standalone)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=10

[Service]
Type=simple
# EnvironmentFile carries DISCORD_TOKEN / DISCORD_APP_ID / BOT_SECRET (required),
# plus DEVNET_URL, DATABASE_URL, DESCENT_ANNOUNCE_CHANNEL_ID, DRAND_API_BASE.
# Placed by ember, chmod 600. Same file as the web unit; the web-only vars are ignored.
EnvironmentFile=%h/.config/dregg/games.env
Environment=RUST_LOG=info
# discord-bot is a STANDALONE cargo workspace (EXCLUDED from the root workspace —
# sqlx/libsqlite3-sys links-conflict; see discord-bot/Cargo.toml + root Cargo.toml
# exclude=). It therefore builds in its OWN target/ dir, NOT the root target/.
WorkingDirectory=%h/dev/breadstuffs/discord-bot
ExecStart=%h/dev/breadstuffs/discord-bot/target/release/dregg-discord-bot
Restart=always
RestartSec=5

# Security hardening (user-unit-safe subset).
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=%h/.local/state/dregg-games
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits — a thin Discord client; bound well below hbox capacity.
LimitNOFILE=65536
MemoryHigh=768M
MemoryMax=1G
TasksMax=512

[Install]
WantedBy=default.target
