# deploy/games/dregg-web-games-funnel.service
# ✅ THIS IS THE UNIT THAT ACTUALLY RUNS. Verified live 2026-07-15.
#   Installed on hbox as a systemd USER unit with `loginctl enable-linger`; published
#   by `tailscale funnel` at https://hbox-dregg.skunk-emperor.ts.net. Both the unit and
#   the funnel config SURVIVED a hard reboot cleanly — this is the reference example of
#   deploy/PRACTICES.md §2 (unit + linger = durable).
#
# ⚠ ITS DREGG_NODE_URL TARGET IS DEAD. The 127.0.0.1:8420 node below was hand-run with
#   an ephemeral --data-dir and did NOT survive the reboot; its ledger is permanently
#   gone. Submitted runs cannot anchor until a node is back on :8420 — and the
#   replacement must be a UNIT with a PERSISTENT data dir, not another hand-run process
#   (TODO-1, deploy/README.md). See deploy/games/RUNBOOK-FUNNEL.md's status banner.
#
# ⚠ hbox is ALSO the build/prove box, and build load hard-killed it this session,
#   taking this demo down with it. Moving the demo off hbox (or capping the prover) is
#   TODO-3. The MemoryHigh/MemoryMax below bound THIS unit — they do nothing to bound
#   the prover, which is what actually kills the box.
#
# ---------------------------------------------------------------------------
# dregg games web demo — FUNNEL VARIANT (Tailscale-Funnel public exposure).
#
# Same `dreggnet-web-server` as dregg-web-games.service, but wired for the
# TAILSCALE-FUNNEL topology instead of the AWS-gateway topology:
#
#   • BIND 127.0.0.1:8790  (NOT the tailnet iface). Tailscale Funnel proxies a
#     LOCAL port: `tailscale funnel 8790` serves 127.0.0.1:8790 publicly at
#     https://hbox-dregg.<tailnet>.ts.net with auto-TLS. Funnel reaches into
#     localhost, so the server binds loopback — contrast the gateway variant,
#     which binds the tailscale iface (100.95.240.73:8790) so a SEPARATE gateway
#     Caddy can reverse-proxy over the tailnet. Here there is NO gateway and NO
#     Caddy: Funnel IS the public edge (TLS + the public hostname).
#
#   • DREGG_NODE_URL=http://127.0.0.1:8420 — the hbox solo devnet node (faucet).
#     A submitted Descent run anchors on that node's ledger (resolve_node_target
#     → NodeTarget::Federation). Unset would keep it in-process (NodeTarget::Local).
#     DREGG_NODE_BEARER is supplied via the env file after unlock (only needed if
#     the node has a passphrase; a loopback devnet needs none).
#
# hbox discipline (deploy/hbox/RUNBOOK.md): a USER unit under
# ~/.config/systemd/user/ , `loginctl enable-linger hbox` so it survives logout,
# `%h`-relative paths, WantedBy=default.target.
#
# ⚠ THE PUBLIC-EXPOSURE FLIP IS NOT HERE. This unit only binds loopback. Nothing
# is public until `tailscale funnel 8790` is run — the single ember-gated
# public-exposure decision (deploy/games/RUNBOOK-FUNNEL.md step d).
[Unit]
Description=dregg games web demo (dreggnet-web-server, funnel variant, loopback:8790)
After=network-online.target
Wants=network-online.target
# Restart-storm brake (same discipline as dregg-web-games.service). A hard crash
# loop stops the unit instead of pegging the box. `systemctl --user reset-failed
# dregg-web-games-funnel` re-arms it after a fix.
StartLimitIntervalSec=600
StartLimitBurst=10

[Service]
Type=simple
# Funnel-variant env file (DREGGNET_WEB_BIND, DREGG_NODE_URL, optional
# DREGG_NODE_BEARER, DATABASE_URL). SEPARATE from the gateway variant's games.env
# so the two topologies never share a bind line. Placed by ember, chmod 600.
EnvironmentFile=%h/.config/dregg/games-funnel.env
Environment=RUST_LOG=info
# Fallback bind if the env file omits DREGGNET_WEB_BIND: LOOPBACK (Funnel proxies
# the local port; the server never binds a public iface). NEVER 0.0.0.0.
Environment=DREGGNET_WEB_BIND=127.0.0.1:8790
# Fallback node target if the env file omits DREGG_NODE_URL: the hbox solo devnet
# node on loopback :8420 (the faucet node). Submitted runs anchor on its ledger.
Environment=DREGG_NODE_URL=http://127.0.0.1:8420
# Durable game sessions: each session's move-log persists here and is resumed on
# boot by replay (a tampered log refuses to reopen). Without this, a restart drops
# every in-progress game.
Environment=DREGGNET_WEB_SESSION_DIR=%h/dregg-games/sessions
# SHARED interaction-envelope audit dir (dreggnet-audit). Set the SAME dir on ALL
# THREE services (this funnel web + dregg-telegram-bot + the discord unit) so one
# `auditq correlate` joins a chain across them — without it each service writes to
# its own sibling `audit/` and cross-service correlate is dead on arrival. Files
# are per-process (audit-DATE.<platform>-<pid>.NN.jsonl), so the shared dir never
# contends. The games-funnel.env file may override; =off disables.
Environment=DREGG_AUDIT_DIR=%h/.local/state/dregg/audit
WorkingDirectory=%h/dev/breadstuffs
ExecStart=%h/dev/breadstuffs/target/release/dreggnet-web-server
Restart=always
RestartSec=5

# Security hardening — DISABLED on hbox, honestly (discovered 2026-07-17 via the
# kubo unit): hbox's user manager cannot apply this sandboxing set ("Failed to
# drop capabilities ... status=218/CAPABILITIES"), and the INSTALLED copy of this
# very unit runs with the whole block commented out — this repo copy previously
# showed the aspirational set, which is exactly the repo-vs-box drift we hunt.
# Keep NoNewPrivileges (works everywhere); re-enable the rest only after
# observing them start ON THE BOX. Note: without ProtectSystem=strict the
# ReadWritePaths lines are inert; the session store needs only
# `mkdir -p ~/dregg-games/sessions` (DREGGNET_WEB_SESSION_DIR).
NoNewPrivileges=true
#hbox ProtectSystem=strict
#hbox ReadWritePaths=%h/.local/state/dregg-games
#hbox ReadWritePaths=-%h/dregg-games
#hbox PrivateTmp=true
#hbox ProtectKernelTunables=true
#hbox ProtectKernelModules=true
#hbox ProtectControlGroups=true

# Resource limits — the games server is a thin in-process web app; bound it well
# below hbox's capacity so it can never starve the co-resident :8420 node / builds.
LimitNOFILE=65536
MemoryHigh=1G
MemoryMax=2G
TasksMax=1024

[Install]
WantedBy=default.target
