# dreggnet/deploy/node-agent/node-agent.service
[Unit]
Description=DreggNet compute-backend bridge agent (lease -> durable metered workflow) on :8021
Documentation=https://ember.software
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dregg
# The native release binary built on the compute node (cargo build --release -p dreggnet-node-agent).
WorkingDirectory=/opt/dreggnet
ExecStart=/opt/dreggnet/target/release/node-agent
Environment=NODE_AGENT_BIND=0.0.0.0:8021
# the native python tier (a seam today) would spawn a python3 subprocess; keep PATH sane.
Environment=PATH=/home/dregg/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Restart=on-failure
RestartSec=3
# Survive reboots; restart if it ever exits.
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
