# dreggnet/deploy/persvati-agent/persvati-agent.service
[Unit]
Description=DreggNet compute-backend bridge agent (lease -> durable metered workload) on :8021
Documentation=https://ember.software
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ember
# The native release binary built on persvati (cargo build --release -p dreggnet-persvati-agent).
WorkingDirectory=/home/ember/dev/DreggNet
ExecStart=/home/ember/dev/DreggNet/target/release/persvati-agent
Environment=PERSVATI_AGENT_BIND=0.0.0.0:8021
# Keep PATH sane for any subprocess a workload tier may spawn.
Environment=PATH=/home/ember/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Restart=on-failure
RestartSec=3
# Survive reboots; restart if it ever exits.
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
