# dreggnet/deploy/observability/node-thermal-exporter/node-thermal-exporter.service
[Unit]
Description=node thermal Prometheus exporter (CPU temp/freq/fan/load on overlay :8022)
Documentation=https://ember.software
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dregg
ExecStart=/usr/bin/python3 /opt/dreggnet/deploy/observability/node-thermal-exporter/node-thermal-exporter.py
Environment=NODE_THERMAL_HOST=0.0.0.0
Environment=NODE_THERMAL_PORT=8022
Environment=NODE_THERMAL_CSV=/var/log/node-thermal.csv
Restart=on-failure
RestartSec=3
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
