# dreggnet/deploy/node-tuning/node-thermal-config.service
[Unit]
Description=node-a steady-state thermal/CPU config (fan level 3 + CPU caps)
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/node-thermal-config.sh

[Install]
WantedBy=multi-user.target
