# dreggnet/deploy/observability/node-exporter/node_exporter.service
[Unit]
# prometheus node_exporter — host metrics for a DreggNet compute box (node-a,
# node-b). Serves /metrics on 0.0.0.0:9100 so the edge Prometheus can scrape it
# over the headscale overlay (100.64.0.x:9100). The headscale ACL opens 9100 on
# tag:compute from tag:edge (deploy/staging/headscale/acls.hujson). This is the
# host twin of the node-thermal-exporter sidecar — node_exporter gives the
# full CPU/RAM/disk/network/load/uptime surface; the thermal exporter adds the
# box-specific temp/freq/fan series node_exporter cannot read.
Description=Prometheus node_exporter (DreggNet host metrics)
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Unprivileged, no shell, no home — node_exporter only reads /proc and /sys.
User=node_exporter
Group=node_exporter
ExecStart=/usr/local/bin/node_exporter \
  --web.listen-address=0.0.0.0:9100 \
  --collector.filesystem.mount-points-exclude=^/(sys|proc|dev|run/credentials)($|/) \
  --collector.systemd
Restart=on-failure
RestartSec=5
# Hardening — read-only system view is all it needs.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
PrivateTmp=true

[Install]
WantedBy=multi-user.target
